# load_lxc_db8_on_boot/lxc-start.service
[Unit]
Requires=ls-hubd.service
After=ls-hubd.service

[Service]
Type=simple
ExecStart=/usr/bin/lxc-start -F -n db8_container
